# .config/supervisord.conf
[program:django]
command=python /srv/backend/app/manage.py runserver 0:8000